# src/meta-user/recipes-apps/distillation-project/embedded/embedded.service.in
[Unit]
Description=Embedded application
DefaultDependencies=no
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
ExecStart=/usr/bin/embedded -port 50001 -config @CONFIG@

[Install]
WantedBy=basic.target
