# src/meta-user/recipes-apps/distillation-project/gui/xserver-nodm.service.in
[Unit]
Description=Xserver startup without a display manager
After=distillation.service osservice.service

[Service]
EnvironmentFile=/etc/default/xserver-nodm
User=@USER@
ExecStart=/etc/xserver-nodm/Xserver
Restart=always

[Install]
Alias=display-manager.service
WantedBy=basic.target
